# rootfs/etc/supervisor/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f -L 15
autostart=true
autorestart=true
priority=30
stdout_logfile=/logs/cron-stdout.log
stderr_logfile=/logs/cron-stderr.log
