# rootfs/etc/supervisor/conf.d/forward-bot.conf
[program:forward-bot]
command=/bin/sh -lc 'if [ -f /config/forward_bot.json ]; then if [ -f /app/forward_bot_runtime.py ]; then exec /usr/bin/python3 /app/forward_bot_runtime.py; elif [ -f /app/forward_bot.py ]; then exec /usr/bin/python3 /app/forward_bot.py; fi; fi; echo "forward bot disabled"; exec sleep infinity'
autostart=true
autorestart=true
startsecs=5
priority=25
stdout_logfile=/logs/forward-bot.log
stderr_logfile=/logs/forward-bot.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=1
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=1
environment=PYTHONUNBUFFERED="1"
