# rootfs/etc/supervisor/conf.d/scanner.conf
[program:scanner-watcher]
command=/app/watcher.sh
autostart=true
autorestart=true
priority=20
stdout_logfile=/logs/scanner.log
stderr_logfile=/logs/scanner.log
environment=PYTHONUNBUFFERED="1"
